# fb_messager/fb-webhooks.service
[Unit]
Description=Facebook & Messenger Webhooks Node.js Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/home/baobao/Projects/Webhooks/fb_messager
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=fb-webhooks
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
